# images/hermes/sandbox-workload.service
[Unit]
Description=OpenShell sandbox workload (Hermes, network-only sibling)
After=openshell-sandbox.service
BindsTo=openshell-sandbox.service
# Set in openshell-sandbox ExecStartPre when SUPERVISOR_MODE=network.
ConditionPathExists=/run/openshell/want-workload

[Service]
Type=simple
ExecStart=/usr/local/lib/openshell/sandbox-workload-run.sh
Restart=on-failure
RestartSec=5

[Install]
# Pulled in when openshell-sandbox starts (not multi-user). After= waits until
# the supervisor unit is active, which is after ExecStartPre wrote the gate.
WantedBy=openshell-sandbox.service
